# ai-che.service
[Unit]
Description=AI Che FastAPI Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/AI-CHE
Environment="PYTHONUNBUFFERED=1"
ExecStart=/root/AI-CHE/venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8000 --workers 4

# ── Restart policy ─────────────────────────────────────────────────────────
# on-failure (а не always) — не перезапускаемся при чистом выходе по SIGTERM.
# RestartSec=10 — при шторме крашей не молотим БД/AI-провайдеры каждые 3 сек.
# StartLimitInterval/Burst — после 5 крашей за минуту systemd сделает паузу.
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5

# ── Resource limits ────────────────────────────────────────────────────────
LimitNOFILE=65536

# ── Logs ───────────────────────────────────────────────────────────────────
StandardOutput=journal
StandardError=journal

# ── Hardening (sandbox) ────────────────────────────────────────────────────
# Отдельный пользователь сильно усложняет post-exploitation. Создание:
#   sudo useradd -r -s /usr/sbin/nologin -d /root/AI-CHE aiche
#   sudo chown -R aiche:aiche /root/AI-CHE
# После этого раскомментировать и сделать `systemctl daemon-reload`.
# (Пока закомментировано чтобы не сломать существующий деплой —
#  включается отдельным шагом миграции.)
# User=aiche
# Group=aiche

# Запрет escalation привилегий (даже если бинарник setuid'ный)
NoNewPrivileges=true

# Изолированный /tmp и /var/tmp на каждый запуск процесса
PrivateTmp=true

# Read-only весь / кроме явно разрешённых путей. ProtectSystem=strict
# монтирует / / /usr / /boot / /etc как ro.
ProtectSystem=strict

# /home, /root и /run/user недоступны для чтения процессу
ProtectHome=true

# Скрыть kernel internals от процесса
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true

# Запрет realtime/SETUID/SYSADMIN/MAC capabilities
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true

# Сетевые семейства: только IPv4/IPv6 (без AF_UNIX, AF_NETLINK и т.д.)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Системные вызовы: только основной набор (отсекает kernel-exploits)
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# Пути к которым процессу нужна запись (БД, аплоады, бэкапы, логи прокси)
ReadWritePaths=/root/AI-CHE/chat.db /root/AI-CHE/chat.db-wal /root/AI-CHE/chat.db-shm
ReadWritePaths=/root/AI-CHE/uploads /root/AI-CHE/backups /root/AI-CHE/proxy_debug.log

# Если включаем User=aiche, эти пути должны быть chown'ены на aiche

[Install]
WantedBy=multi-user.target
