# ai-che.service
[Unit]
Description=AI Che FastAPI Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/AI-CHE
Environment="PYTHONUNBUFFERED=1"
# uvicorn слушает ТОЛЬКО 127.0.0.1: внешние запросы идут через nginx (HTTPS,
# CSP, headers). Раньше был --host 0.0.0.0 → API торчал наружу на :8000
# обходя nginx-уровень защиты.
ExecStart=/root/AI-CHE/venv/bin/python -m uvicorn main:app --host 127.0.0.1 --port 8000 --workers 4

# ── Restart policy ─────────────────────────────────────────────────────────
# on-failure (а не always) — не перезапускаемся при чистом выходе по SIGTERM.
# RestartSec=10 — при шторме крашей не молотим БД/AI-провайдеры каждые 3 сек.
# StartLimitInterval/Burst — после 5 крашей за минуту systemd сделает паузу.
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5

# ── Resource limits ────────────────────────────────────────────────────────
LimitNOFILE=65536

# ── Logs ───────────────────────────────────────────────────────────────────
StandardOutput=journal
StandardError=journal

# ── Hardening (минимальный безопасный набор) ───────────────────────────────
# Полный набор Protect*/Restrict*/SystemCallFilter ломал uvicorn workers
# (см. инцидент 2026-04-27: ProtectHome=true → ExecStart 203/EXEC, после
# фикса — workers crash на startup из-за ProtectSystem=strict +
# MemoryDenyWriteExecute). Включаем только то, что точно не ломает.
#
# Полный hardening — отдельным спринтом миграции на User=aiche и /opt/ai-che,
# с поэтапным включением каждой директивы.

# Запрет escalation привилегий (даже если бинарник setuid'ный)
NoNewPrivileges=true

# Изолированный /tmp и /var/tmp на каждый запуск процесса
PrivateTmp=true

[Install]
WantedBy=multi-user.target
